# dining/var/www/html/ds18b20.service
[Unit]
Description=Read DS18B20 Temperature sensor
[Service]
Type=simple
ExecStart=/var/www/html/ds18b20.sh
