# dining/var/www/html/PindaNetUpdate.service
[Unit]
Description=Update and Reset
[Service]
Type=simple
ExecStart=/var/www/html/PindaNetUpdate.sh
