# domoticaController/var/www/html/checkWiFi.service
[Unit]
Description=Check WiFi connection
[Service]
Type=simple
ExecStart=/var/www/html/checkWiFi.sh
