# etc/systemd/system/PindaNetBluetoothDetection.service
[Unit]
Description=Bluetooth Detection
[Service]
Type=simple
ExecStart=/usr/sbin/bluetooth-detection.sh
