# etc/systemd/system/PindaNetThermostat.service
[Unit]
Description=Thermostat

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/html
ExecStart=/var/www/html/thermostat.sh
