# etc/systemd/system/PindaNetWallpaper.service
[Unit]
Description=Get new wallpaper every day
[Service]
Type=simple
ExecStart=/usr/sbin/wallpaper.sh
