# etc/systemd/system/PindaNetBluetoothScan.service
[Unit]
Description=Bluetooth Detection Scan
[Service]
Type=simple
ExecStart=/usr/sbin/bluetoothscan.sh
