# etc/systemd/system/PindaNetThermostat.service
[Unit]
Description=Thermostat
[Service]
Type=simple
ExecStart=/usr/sbin/PindaNetThermostat.sh
