# etc/systemd/system/PindaNetState.service
[Unit]
Description=Controller state        
[Service]
Type=simple
ExecStart=/usr/sbin/PindaNetState.sh
