# etc/systemd/system/PindaNetMotion.service
[Unit]
Description=Motion Detection

[Service]
ExecStart=/usr/sbin/PindaNetMotion.py

[Install]
WantedBy=multi-user.target
