# etc/systemd/system/PindaNetThermostat.timer
[Unit]
Description=Thermostat
[Timer]
OnBootSec=1min
OnUnitActiveSec=5min
Unit=PindaNetThermostat.service
[Install]
WantedBy=multi-user.target
