# etc/systemd/system/PindaNetMotionDetection.service
[Unit]
Description=Motion Detection
[Service]
Type=simple
ExecStart=/usr/sbin/motion-detection.sh
